# usbguard-notifier.service.in
[Unit]
Description=USBGuard Notifier
After=usbguard.service

[Service]
ExecStart=%bindir%/usbguard-notifier -w

[Install]
WantedBy=default.target
