# usbguard_notifier.service.in
[Unit]
Description=USBGuard Notificator
After=usbguard.service

[Service]
ExecStart=%bindir%/usbguard_notifier

[Install]
WantedBy=default.target
